# claude-telegram-bridge/com.yourname.claude-telegram-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yourname.claude-telegram-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update this path after running setup.sh -->
        <string>/Users/YOUR_USERNAME/claude-telegram-bridge/venv/bin/python</string>
        <string>/Users/YOUR_USERNAME/claude-telegram-bridge/bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/claude-telegram-bridge</string>

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/claude-telegram-bridge/logs/bot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/claude-telegram-bridge/logs/bot.error.log</string>

    <!-- Throttle restarts (seconds) -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
